# supervisord.conf
[supervisord]
nodaemon=true

[program:php-serve]
command=php -S 0.0.0.0:8001 -t public
directory=/var/www/html
autostart=true
autorestart=true


